# recipes-app/alarm/alarm/alarm_snooze.timer
[Unit]
Description=Alarm Snooze Timer

[Timer]
OnCalendar=*-*-* 07:35

[Install]
WantedBy=multi-user.target
