# recipes-app/quetzalcoatl/quetzalcoatl/mpc_mediaserver.service
[Unit]
Description=MPC Qt5 

[Service]
ExecStartPre=/bin/systemctl stop start
ExecStartPre=/bin/sleep 0.1

EnvironmentFile=/etc/environment
ExecStart=/opt/quetzalcoatl -platform eglfs

ExecStartPost=/bin/sleep 0.1
ExecStopPost=/bin/systemctl start start

[Install]
WantedBy=multi-user.target
